# assets/speedtest.timer
[Unit]
Description=Run speedtest timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target